# cbdRegional/Web/[MEDIUM] ponder/dist/conf/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:ponder]
command=gunicorn --bind 0.0.0.0:5000 --workers 2 --worker-class gthread --threads 4 --timeout 120 --graceful-timeout 30 --keep-alive 5 run:app
directory=/opt/webapp
user=appuser
environment=PYTHONPATH="/opt/webapp",HOME="/home/appuser",XDG_CONFIG_HOME="/home/appuser/.config",BOT_URL="%(ENV_BOT_URL)s",SECRET_KEY="%(ENV_SECRET_KEY)s",FLAG="%(ENV_FLAG)s",ADMIN_USERNAME="%(ENV_ADMIN_USERNAME)s",ADMIN_PASSWORD="%(ENV_ADMIN_PASSWORD)s",ADMIN_UPLOAD_FILENAME="%(ENV_ADMIN_UPLOAD_FILENAME)s"
priority=10
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bot]
command=node app.js
directory=/opt/bot
user=appuser
environment=HOME="/home/appuser",XDG_CONFIG_HOME="/home/appuser/.config",PORT="3001",HOST="127.0.0.1",APP_URL="%(ENV_APP_URL)s",ADMIN_USERNAME="%(ENV_ADMIN_USERNAME)s",ADMIN_PASSWORD="%(ENV_ADMIN_PASSWORD)s",PUPPETEER_EXECUTABLE_PATH="%(ENV_PUPPETEER_EXECUTABLE_PATH)s"
autostart=true
autorestart=true
priority=20
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:state-reset]
command=/usr/local/bin/reset-state-loop.sh
user=root
autostart=true
autorestart=true
priority=30
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
